# docker/selfhost/rootfs/etc/s6-overlay/s6-rc.d/minio/run
#!/command/with-contenv bash
# MinIO — S3-compatible storage. Bound to 0.0.0.0 because browsers load avatars /
# screenshots directly from S3_PUBLIC_URL (e.g. http://localhost:9000/intlayer/...).
export MINIO_ROOT_USER="${S3_ACCESS_KEY_ID}"
export MINIO_ROOT_PASSWORD="${S3_SECRET_ACCESS_KEY}"
exec minio server /data/minio \
  --address ":9000" \
  --console-address ":9001"
