# docker/selfhost/rootfs/etc/s6-overlay/s6-rc.d/init-minio/up
/command/with-contenv /etc/s6-overlay/scripts/init-minio.sh
